# testeDockerBD/crontab
* * * * *  echo "ola" >> /var/log/cron.log 
